# StealthShark_NFC_Combined/com.aimf.stealthshark-nfc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aimf.stealthshark-nfc</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/flowgirl/Documents/_MobileShield/StealthShark_NFC_Combined/stealthshark_nfc_combined.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/flowgirl/Documents/_MobileShield/StealthShark_NFC_Combined/logs/stealthshark-nfc.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/flowgirl/Documents/_MobileShield/StealthShark_NFC_Combined/logs/stealthshark-nfc.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/flowgirl/Documents/_MobileShield/StealthShark_NFC_Combined</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/flowgirl/Documents/_MobileShield/StealthShark_NFC_Combined</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>
</dict>
</plist>
