# Desktop_Shortcuts/com.aimfllc.stealthshark.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aimfllc.stealthshark</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/flowgirl/Documents/StealthShark-PublicRelease/launch.sh</string>
        <string>--auto-start</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/flowgirl/Documents/StealthShark-PublicRelease</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/stealthshark.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/stealthshark_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/flowgirl/Documents/StealthShark-PublicRelease</string>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
