# backups/backup_20250831_151610/com.stealthshark.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stealthshark.monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/flowgirl/Documents/StealthShark && source venv/bin/activate 2>/dev/null || python3 -m venv venv && source venv/bin/activate && pip3 install -r requirements.txt; python3 enhanced_memory_monitor.py --duration 4</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/flowgirl/Documents/StealthShark/logs/stealthshark.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/flowgirl/Documents/StealthShark/logs/stealthshark_error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/flowgirl/Documents/StealthShark</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/flowgirl/Documents/StealthShark</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
